# lab-7/crontab
* * * * * Information-Security-2025/lab-7/mail.py
